# root/etc/s6-overlay/s6-rc.d/acquire_airconnect/up
/etc/s6-overlay/scripts/acquire_airconnect_up
