# root/etc/s6-overlay/s6-rc.d/airupnp/run
#!/bin/sh

echo "Starting airupnp run"


if [ -z "${AIRUPNP_VAR}" ]; then
   /bin/airupnp-"${ARCH_VAR}" -l 1000:2000
else
   /bin/airupnp-"${ARCH_VAR}" "${AIRUPNP_VAR}"
fi


