# root/etc/s6-overlay/s6-rc.d/airupnp/run
#!/usr/bin/with-contenv bash

echo "Starting airupnp run"

if [ "$ARCH_VAR" == "amd64" ]; then
  ARCH_VAR=linux-x86_64
elif [ "$ARCH_VAR" == "arm64" ]; then
  ARCH_VAR=linux-aarch64
# elif [ "$ARCH_VAR" == "arm" ]; then
#   ARCH_VAR=linux-arm
fi

echo "Checking for valid arch options"
case $ARCH_VAR in
  linux-x86_64)
    echo "Proceeding with linux-x86_64 arch"
    ;;
  linux-aarch64)
    echo "Proceeding with linux-aarch64 arch"
    ;;
  # linux-arm)
  #   echo "Proceeding with linux-arm arch"
  #   ;;
  *)
    echo "Unrecognized or invalid arch selection, CANCELING INSTALL"
    echo "========== FAILURE DETECTED ========="
    echo "YOUR CONTAINER WILL NOT WORK, PLEASE ADDRESS OR OPEN AN ISSUE"
    exit 1
    ;;
esac

echo "ARCH_VAR=${ARCH_VAR}"
ARCH_VAR="${ARCH_VAR}-static"
echo "ARCH_VAR=${ARCH_VAR}"

if [ -f /bin/airupnp-"${ARCH_VAR}" ]; then    
    echo "executable /bin/airupnp-"${ARCH_VAR}" exists"
    if [ -z "${AIRUPNP_VAR}" ]; then
        echo "AIRUPNP_VAR empty, setting it to defaults"
        AIRUPNP_VAR="-l 1000:2000"
        echo "AIRUPNP_VAR=${AIRUPNP_VAR}"
    fi
    echo "/bin/airupnp-${ARCH_VAR} -h"
    /bin/airupnp-${ARCH_VAR} -h
    echo "/bin/airupnp-${ARCH_VAR} ${AIRUPNP_VAR}"
    /bin/airupnp-${ARCH_VAR} ${AIRUPNP_VAR}
fi

echo "Exiting airupnp run"
