# root/etc/s6-overlay/s6-rc.d/airupnp/run
#!/bin/sh

echo "Starting airupnp run"

/bin/airupnp-"${ARCH_VAR}" -l 1000:2000
/bin/airupnp-"${ARCH_VAR}" "${AIRUPNP_VAR}"
