# root/etc/s6-overlay/s6-rc.d/aircast/run
#!/usr/bin/with-contenv bash

echo "Starting aircast run"

if [ ! -f /bin/aircast-"${ARCH_VAR}" ]; then
    echo "executable not found - exiting"
    exit 1
else
    echo "file exists"
    if [ -z "${AIRCAST_VAR}" ]; then
       /bin/aircast-"${ARCH_VAR}"
    else
       /bin/aircast-"${ARCH_VAR}" "${AIRCAST_VAR}"
    fi
fi

echo "Exiting aircast run"
